# smartwatch/etc/systemd/system/smartwatch.timer
[Unit]
Description=smartwatch

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
